# quadlet/dev/learn-language-db.container
[Unit]
Description=PostgreSQL development database

[Container]
ContainerName=learn-language-db
Image=docker.io/library/postgres:17.5-bullseye
Environment=POSTGRES_DB=learnlanguage
Environment=POSTGRES_USER=postgres
Environment=POSTGRES_PASSWORD=postgres
PublishPort=5433:5432
Volume=learn-language-postgres-data:/var/lib/postgresql/data
Volume={{PROJECT_DIR}}/server/src/main/resources/schema.sql:/docker-entrypoint-initdb.d/init.sql
HealthCmd=pg_isready -U postgres
HealthInterval=10s
HealthTimeout=5s
HealthRetries=5

[Install]
WantedBy=default.target
