# quadlet/test/mock-anthropic.container
[Unit]
Description=Mock Anthropic server

[Container]
ContainerName=mock-anthropic
Image=localhost/learn-language-mock-anthropic:test
Pod=learn-language-test.pod
Environment=PORT=3003
HealthCmd=curl --fail --silent localhost:3003/health || exit 1
HealthInterval=2s
HealthTimeout=3s
HealthRetries=5
HealthStartPeriod=2s

[Install]
WantedBy=default.target
